# src/com.user.stockmarketagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.stockmarketagent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/goyamegh/workplace/stock_market_agent/src/run_stock_agent.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/goyamegh/workplace/stock_market_agent</string>
    <key>StandardOutPath</key>
    <string>/Users/goyamegh/workplace/stock_market_agent/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/goyamegh/workplace/stock_market_agent/logs/stderr.log</string>
</dict>
</plist> 